# com.example.claude-usage-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  LaunchAgent template. Before use:
    1. Replace __HOME__ with your home path (e.g. /Users/alice).
    2. Replace __PYTHON__ with `which python3` output (e.g. /opt/homebrew/opt/python@3.13/bin/python3 on Apple Silicon with Homebrew).
    3. Save as ~/Library/LaunchAgents/com.<you>.claude-usage-report.plist and update the <Label> accordingly.
    4. Load: launchctl load ~/Library/LaunchAgents/com.<you>.claude-usage-report.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.claude-usage-report</string>

  <key>ProgramArguments</key>
  <array>
    <string>__PYTHON__</string>
    <string>__HOME__/.claude/scripts/claude-usage-report.py</string>
    <string>--refresh</string>
    <string>30</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>__HOME__</string>
  </dict>

  <key>StartInterval</key>
  <integer>30</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/claude-usage-report.log</string>

  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/claude-usage-report.log</string>

  <key>WorkingDirectory</key>
  <string>__HOME__</string>
</dict>
</plist>
